# serviço/adic.service
[Unit]
Description=Adicionar usuarios

[Service]
ExecStart=/usr/sbin/adicionar
Type=simple

[Install]
WantedBy=default.target
